# systemd-user/end-blue-light-filter.timer
[Unit]
Description=End Blue Light Filter

[Timer]
OnCalendar=*-*-* 14:00:00
Persistent=true

[Install]
WantedBy=target.target
