# systemd-user/start-blue-light-filter.timer
[Unit]
Description=Start Blue Light Filter

[Timer]
OnCalendar=*-*-* 5:00:00
Persistent=true

[Install]
WantedBy=target.target
