# systemd-user/waydroid-check.timer
[Unit]
Description=Run Waydroid check every minute

[Timer]
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
