# systemd-user/waydroid-check.timer
[Unit]
Description=Run Waydroid check every 30 minutes

[Timer]
OnUnitInactiveSec=1min

[Install]
WantedBy=timers.target
